# Assignments/week06/remote-job-scrapper/remote-job-scraper-api-deb/DEBIAN/prerm
#!/bin/bash
set -e

echo "Stopping remote-job-scraper-api service..."
systemctl stop remote-job-scraper-api.service || true
systemctl disable remote-job-scraper-api.service || true

echo "Removing ufw rule for port 8889..."
if command -v ufw > /dev/null; then
    ufw delete allow 8889/tcp || true
    echo "Port 8889/tcp removed from ufw."
fi

exit 0
